# deploy/systemd/traderie-ingest-snapshot.timer
# =============================================================================
# traderie-ingest-snapshot.timer
# DO NOT INSTALL OR ENABLE WITHOUT SCHEDULER GATE
# =============================================================================
# Production cadence: 4x daily.
# Install ONLY after SHARED-004 capacity preflight and Scheduler Gate approval.

[Unit]
Description=Traderie ingest snapshot timer (4x daily)
Documentation=https://github.com/wguDataNinja/traderie

[Timer]
OnCalendar=00:00:00
OnCalendar=06:00:00
OnCalendar=12:00:00
OnCalendar=18:00:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
