# deploy/systemd/traderie-ingest-snapshot.service
# =============================================================================
# traderie-ingest-snapshot.service
# DO NOT INSTALL OR ENABLE WITHOUT SCHEDULER GATE
# =============================================================================
# Runs a full generation across all 4 segments sequentially.
# Each segment is independently bounded by a per-segment timeout in
# scripts/run_traderie_generation.sh (pc_sc_nl=180s, pc_sc_l=240s,
# pc_hc_l=360s, pc_hc_nl=480s). The generation service itself has no
# timeout because individual segments control their own runtime.

[Unit]
Description=Traderie ingest — full generation across all 4 segments
After=network-online.target
Documentation=https://github.com/wguDataNinja/traderie

[Service]
Type=oneshot
User=scraper
WorkingDirectory=/home/scraper/apps/traderie
EnvironmentFile=/home/scraper/config/traderie.env

# Flock prevents concurrent generation runs
ExecStartPre=/usr/bin/mkdir -p /home/scraper/data/traderie/.locks
ExecStart=/usr/bin/flock -n /home/scraper/data/traderie/.locks/snapshot.lock \
    /home/scraper/apps/traderie/scripts/run_traderie_generation.sh

# Each segment has its own timeout in the orchestrator script.
# No generation-level timeout needed — individual segments bound themselves.
TimeoutStartSec=infinity

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=timers.target
